# ocs_ci/utility/networksplit/systemd/network-split-teardown@.timer
[Unit]
Description=Network split teardown timer

[Timer]
OnCalendar=@%i
Unit=network-split-teardown.service

[Install]
WantedBy=timers.target
